# docs/policy-app.service
# systemd 서비스 예시 (AI 공약 멘토링)
# 복사: sudo cp policy-app.service /etc/systemd/system/
# 수정: User, WorkingDirectory, ExecStart 경로를 본인 서버에 맞게 변경
# 적용: sudo systemctl daemon-reload && sudo systemctl enable policy-app && sudo systemctl start policy-app

[Unit]
Description=AI 공약 멘토링 서버
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/Policy
Environment="PATH=/home/ec2-user/Policy/.venv/bin"
ExecStart=/home/ec2-user/Policy/.venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000 --timeout-keep-alive 300
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
